# data/rc.local
#!/bin/bash

/root/scripts/10-dhcp-client.sh
/root/scripts/20-setup-nvidia.sh
/root/scripts/20-setup-amd.sh
/root/scripts/20-setup-sensors.sh

exit 0
